# tool/pi_setup/rc.local

(cd /home/pi/work/udpnode && ./run_deploy.sh &)
