# tool/pi_setup/rc.local

(cd /home/pi/work/bundle && ./run_deploy.sh &)
